# diody/crontabs/rainbow_wi-fi_clients
MAILTO=""
*	*	*	*	*	root	/root/diody/wi-fi_clients.sh

