# diody/crontabs/rainbow_cuckoo
MAILTO=""
0	*	*	*	*	root	/root/diody/cron-cuckoo.sh

