# diody/crontabs/rainbow_wi-fi_clients
MAILTO=""
0-59/2	*	*	*	*	root	/root/diody/wi-fi_clients.sh
1-59/2	*	*	*	*	root	/etc/init.d/rainbow restart

