# diody/crontabs/rainbow_wi-fi_clients
MAILTO=""
2-59   *   *   *   *   root    /root/diody/wi-fi_clients.sh
